# rusnas-storage-analyzer/rusnas-storage-collector.timer
[Unit]
Description=rusNAS Storage Analyzer — hourly collector timer

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
