# rusnas-metrics/rusnas-metrics.service
[Unit]
Description=rusNAS Metrics Server (port 9100)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/lib/rusnas
ExecStart=/usr/local/lib/rusnas/metrics_server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
