# rusnas-guard/service/rusnas-guard.service
[Unit]
Description=rusNAS Guard - Ransomware Protection Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/usr/lib/rusnas-guard
ExecStart=/usr/lib/rusnas-guard/guard
Restart=on-failure
RestartSec=5
User=root
RuntimeDirectory=rusnas-guard
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
