# rusnas-guard/service/rusnas-guard.service
[Unit]
Description=rusNAS Guard - Ransomware Protection Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/lib/rusnas-guard/guard.py
Restart=on-failure
RestartSec=5
User=root
RuntimeDirectory=rusnas-guard
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
